# scripts/krb5-renew.timer
[Unit]
Description=Renew Kerberos ticket every 8 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
